# launchd/com.shuta.lex.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shuta.lex</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/shuta/lex-start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/shuta</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/shuta</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- ログイン時に自動起動 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- クラッシュ時に自動再起動（正常終了時は再起動しない） -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- 再起動間隔（秒）- 連続クラッシュ防止 -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- ログ出力先 -->
    <key>StandardOutPath</key>
    <string>/tmp/lex-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/lex-stderr.log</string>
</dict>
</plist>
